# systemd/openbrain-papers-to-prompts.service
[Unit]
Description=openbrain: promote top-scored papers into the morning queue
After=network.target openbrain-research-score.service
Documentation=https://github.com/cgallic/openbrain

[Service]
Type=oneshot
EnvironmentFile=-/etc/openbrain/openbrain.env
ExecStart=/usr/bin/env python3 -m openbrain.scanners.papers_to_prompts
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120
