# systemd/openbrain-today-briefing.timer
[Unit]
Description=Run today-briefing scanner before the morning digest

[Timer]
OnCalendar=*-*-* 11:10:00 UTC
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
