# systemd/openbrain-research-ingest.service
[Unit]
Description=openbrain: ingest arXiv + HF Daily + Semantic Scholar papers
After=network.target
Documentation=https://github.com/cgallic/openbrain

[Service]
Type=oneshot
EnvironmentFile=-/etc/openbrain/openbrain.env
ExecStart=/usr/bin/env python3 -m openbrain.ingest.research_papers --source all --verbose
StandardOutput=journal
StandardError=journal
TimeoutStartSec=900
