# systemd/openbrain-prompt-digest.service
[Unit]
Description=openbrain: daily morning Discord digest (news + papers + topical)
After=network.target openbrain-papers-to-prompts.service openbrain-news-opinion.service openbrain-topical-questions.service openbrain-today-briefing.service
Documentation=https://github.com/cgallic/openbrain

[Service]
Type=oneshot
EnvironmentFile=-/etc/openbrain/openbrain.env
ExecStart=/usr/bin/env python3 -m openbrain.digest.prompt_digest
StandardOutput=journal
StandardError=journal
TimeoutStartSec=180
