# systemd/openbrain-news-opinion.service
[Unit]
Description=openbrain: RSS news scanner -> hot-take prompts
After=network.target ollama.service
Documentation=https://github.com/cgallic/openbrain

[Service]
Type=oneshot
EnvironmentFile=-/etc/openbrain/openbrain.env
ExecStart=/usr/bin/env python3 -m openbrain.scanners.news_opinion
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
