# systemd/openbrain-research-score.service
[Unit]
Description=openbrain: score ingested papers against research-interests.md
After=network.target ollama.service openbrain-research-ingest.service
Documentation=https://github.com/cgallic/openbrain

[Service]
Type=oneshot
EnvironmentFile=-/etc/openbrain/openbrain.env
ExecStart=/usr/bin/env python3 -m openbrain.score.research_papers --verbose
StandardOutput=journal
StandardError=journal
TimeoutStartSec=7200
