# systemd/paperboy-research-digest.timer
[Unit]
Description=Daily top-papers Discord post at 12:00 UTC

[Timer]
OnCalendar=*-*-* 12:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
