# systemd/paperboy-news-opinion.service
[Unit]
Description=paperboy: RSS news scanner -> hot-take prompts
After=network.target ollama.service
Documentation=https://github.com/cgallic/paperboy

[Service]
Type=oneshot
EnvironmentFile=-/etc/paperboy/paperboy.env
ExecStart=/usr/bin/env python3 -m paperboy.scanners.news_opinion
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600
