# systemd/paperboy-topical-questions.timer
[Unit]
Description=Run topical-questions scanner before the morning digest

[Timer]
OnCalendar=*-*-* 11:08:00 UTC
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
