# systemd/paperboy-research-score.timer
[Unit]
Description=Daily research-papers scoring at 10:45 UTC

[Timer]
OnCalendar=*-*-* 10:45:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
