# systemd/paperboy-news-opinion.timer
[Unit]
Description=Run paperboy news-opinion scanner twice daily

[Timer]
OnCalendar=*-*-* 11:00:00 UTC
OnCalendar=*-*-* 21:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
